# examples/matrix-irc-nick.timer
[Unit]
Description=Rename IRC nicks periodically

[Timer]
OnCalendar=*-*-* 04:56:20
Persistent=true

[Install]
WantedBy=timers.target

